# maintenance/ppa-cleaner-package/ppa-cleaner.timer
[Unit]
Description=Weekly Ubuntu PPA cleanup check

[Timer]
OnCalendar=Sun *-*-* 04:15:00
RandomizedDelaySec=30m
Persistent=true
Unit=ppa-cleaner.service

[Install]
WantedBy=timers.target
